# deploy/helpbot.service
[Unit]
Description=HelpBot — Telegram Bot + Admin API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/encode/embassybot
EnvironmentFile=/home/encode/embassybot/.env
ExecStart=/home/encode/embassybot/venv/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
